# crontab
0 0 * * * root . /root/project_env.sh; /backup.sh >> /var/log/cron.log 2>&1
