# crontab
0 */12 * * * root . /root/project_env.sh; /backup.sh >> /var/log/cron.log 2>&1

